# deploy/systemd/intel-dark-ingest.service
[Unit]
Description=BorealSec Intel — ingest_dark (oneshot, Tor)
After=network-online.target intel-tor.service
Wants=network-online.target intel-tor.service

[Service]
Type=oneshot
User=appuser
WorkingDirectory=/opt/intel/borealsec-intel
EnvironmentFile=/opt/intel/.env
ExecStart=/opt/intel/venv/bin/python manage.py ingest_dark \
    --settings=config.settings.prod
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
