# deploy/tor/intel-tor.service
[Unit]
Description=Intel Tor Proxy Container
After=network-online.target
Wants=network-online.target

[Service]
User=appuser
Restart=always
RestartSec=10
ExecStartPre=-/usr/bin/podman rm -f intel-tor
ExecStart=/usr/bin/podman run \
    --rm \
    --name intel-tor \
    --network slirp4netns \
    --publish 127.0.0.1:9050:9050 \
    --cap-drop ALL \
    --security-opt no-new-privileges \
    --read-only \
    --tmpfs /tmp \
    localhost/intel-tor:latest
ExecStop=/usr/bin/podman stop intel-tor

[Install]
WantedBy=multi-user.target
